# docker/supervisor/supervisord.conf
; ============================================================================
; WEBGRE3 - Supervisor Configuration
; Manages PHP-FPM and Nginx processes
; ============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ============================================================================
; PHP-FPM Process
; ============================================================================
[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize --fpm-config /usr/local/etc/php-fpm.conf
autostart=true
autorestart=true
priority=5
stdout_logfile=/var/www/html/storage/logs/php-fpm.stdout.log
stderr_logfile=/var/www/html/storage/logs/php-fpm.stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
user=root

; ============================================================================
; Nginx Process
; ============================================================================
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/nginx/stdout.log
stderr_logfile=/var/log/nginx/stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
user=root

; ============================================================================
; Cron Job Runner (opzionale)
; ============================================================================
[program:cron]
command=/usr/sbin/crond -f -l 8
autostart=true
autorestart=true
priority=15
stdout_logfile=/var/www/html/storage/logs/cron.stdout.log
stderr_logfile=/var/www/html/storage/logs/cron.stderr.log
user=root

; ============================================================================
; Laravel Queue Worker (opzionale, se usi queue)
; ============================================================================
; [program:queue-worker]
; command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600
; autostart=true
; autorestart=true
; priority=20
; numprocs=1
; user=www-data
; stdout_logfile=/var/www/html/storage/logs/queue-worker.log
; stderr_logfile=/var/www/html/storage/logs/queue-worker.error.log

; ============================================================================
; Process Groups
; ============================================================================
[group:webgre3]
programs=php-fpm,nginx,cron
priority=999
